# confs/opendevdata-worker.conf
export PATH="/root/.rbenv/bin:$PATH"
start on starting opendevdata
stop on stopping opendevdata
