# confs/opendevdata-worker-3.conf
start on starting opendevdata-worker
stop on stopping opendevdata-worker
respawn

chdir /var/www/opendevdata

script
    exec /root/.rbenv/shims/rake RAILS_ENV=production environment resque:work QUEUE="*" >> /var/log/opendevdata/worker-3.log 2>&1
end script
